# install/systemd/fs42-remote-controller.service.template
[Unit]
Description=FieldStation42 Remote Controller
After=network.target

[Service]
Type=simple
WorkingDirectory=__INSTALL_DIR__
ExecStart=__INSTALL_DIR__/env/bin/python3 fs42/pi/remote_controller.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
